# init/remtroll_systemd.service
[Unit]
Description=RemTroll

[Service]
ExecStart=/usr/bin/remtroll

[Install]
Requires=network.target
After=network.target
